# docs/launchd-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aluarius.memory-compiler-maintenance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/Users/exmac/Desktop/ideas/claude-memory-compiler</string>
        <string>python</string>
        <string>scripts/maintenance.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/exmac/Desktop/ideas/claude-memory-compiler</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/memory-compiler-maintenance.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/memory-compiler-maintenance.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
